# Mac/Library/LaunchAgents/bitwarden-ssh-env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>bitwarden-ssh-env</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>launchctl setenv SSH_AUTH_SOCK "/Users/federico/Library/Containers/com.bitwarden.desktop/Data/.bitwarden-ssh-agent.sock"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
